# msys/1.0/etc/profile.d/bashrc.default
#!/bin/sh
## Alias python in interactive mode
alias py='py -i'
